# files/default/newrelic-sysmond-rhel.conf
description "New Relic Server Monitor"

start on [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 20

exec sudo -u newrelic /usr/sbin/nrsysmond -f -c /etc/newrelic/nrsysmond.cfg
